# addons/music-assistant-skill/rootfs/etc/s6-overlay/s6-rc.d/init-alexa-skill/run
#!/usr/bin/with-contenv bashio
bashio::log.info "init-alexa-skill: setting up persistent ASK CLI credential store..."

SKILL_HOSTNAME="$(bashio::config 'SKILL_HOSTNAME' || true)"
if bashio::var.is_empty "${SKILL_HOSTNAME}"; then
    bashio::log.fatal "SKILL_HOSTNAME is required. Set it in the add-on configuration and restart."
    exit 1
fi

mkdir -p /data/.ask

if [ -L /root/.ask ]; then
    bashio::log.debug "init-alexa-skill: /root/.ask symlink already present, skipping."
elif [ -d /root/.ask ]; then
    bashio::log.info "init-alexa-skill: migrating existing /root/.ask into /data/.ask..."
    cp -a /root/.ask/. /data/.ask/ 2>/dev/null || true
    rm -rf /root/.ask
    ln -s /data/.ask /root/.ask
else
    ln -s /data/.ask /root/.ask
fi

bashio::log.info "init-alexa-skill: /root/.ask -> /data/.ask (credentials persist across restarts)"
bashio::log.info "init-alexa-skill: done."
