# addons/music-assistant-skill/rootfs/etc/s6-overlay/s6-rc.d/init-alexa-skill/up
/etc/s6-overlay/s6-rc.d/init-alexa-skill/run