# addons/music-assistant-skill/rootfs/etc/s6-overlay/s6-rc.d/alexa-skill/finish
#!/usr/bin/with-contenv bashio
EXIT_CODE="${1}"
SIGNAL="${2}"

if [ "${EXIT_CODE}" -eq 256 ]; then
    bashio::log.error "alexa-skill: killed by signal ${SIGNAL}"
    echo "$((128 + SIGNAL))" > /run/s6-linux-init-container-results/exitcode
else
    if [ "${EXIT_CODE}" -ne 0 ]; then
        bashio::log.error "alexa-skill: exited with non-zero code ${EXIT_CODE}"
        echo "${EXIT_CODE}" > /run/s6-linux-init-container-results/exitcode
    else
        bashio::log.info "alexa-skill: exited cleanly (code 0)"
    fi
fi
