# addons/music-assistant-skill/rootfs/etc/s6-overlay/s6-rc.d/alexa-skill/run
#!/usr/bin/with-contenv bashio
bashio::log.info "alexa-skill: starting Flask application..."

PORT="$(bashio::config 'PORT' '5000' || echo '5000')"

bashio::log.info "alexa-skill: PORT=${PORT}"

exec /app/venv/bin/python /app/app/app.py
